# scripts/smoke-test/com.tfb.smoketest.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tfb.smoketest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/fin/Documents/The Franchisor Blueprint Website/scripts/smoke-test/run-daily.sh</string>
  </array>

  <!--
    Daily 7:00 AM local time. macOS launchd uses the system timezone
    for StartCalendarInterval, so this is automatically 7am whatever
    TZ Eric's Mac is set to (currently America/Denver).
    If the Mac is asleep at 7am, the job runs at the next wake.
  -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/fin/Documents/The Franchisor Blueprint Website</string>

  <!--
    launchd starts with a minimal PATH; the orchestrator script
    re-exports the homebrew + system paths near the top, but we
    also set them here for safety.
  -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/fin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/tfb-smoketest-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tfb-smoketest-stderr.log</string>

  <!--
    Don't fire on load. The job runs only on the cron schedule.
    For an immediate test run after install, use:
      launchctl start com.tfb.smoketest
  -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
